# dot_config/systemd/user/mega-cmd-server.service
[Unit]
Description=MEGA CMD Server Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Systemd user services already know your $USER, so no hardcoded names needed!
ExecStart=/usr/bin/mega-cmd-server
# Automatically restart if the server dies
Restart=on-failure
RestartSec=5

[Install]
# This makes it start as soon as you log in (TTY or Hyprland)
WantedBy=default.target
